# deploy/resume.service
[Unit]
Description=Resume web app (Bun)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/resume
ExecStart=/home/hassan/.bun/bin/bun src/server.ts
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
